# 0_RootFS/Rootfs/bundled/conf/profile.d/aliases.sh
# Enable colors by most of these things by default
alias grep='grep --color=auto'
alias ls='ls --color=auto'
alias ll='ls --color=auto -laF'

# Whenever we say `vi`, we really mean `vim`
alias vi='vim'
