# 0_RootFS/Rootfs/bundled/conf/profile.d/ps1.sh
export PS1="${HIDDEN_PS1:-sandbox \$ }"
